# templates/telescoped.service.template
# This is a template for the systemd .service file.
# Comments will be stripped off when installed.
[Unit]
Description=The Telescope Service for CTMO
After=network.target

[Service]
Type=simple
User=root
# WorkingDirectory=/home/nanodano
ExecStart=scriptname
Restart=on-failure

# [Install]
# WantedBy=otherctmoservice
